# deploy/systemd/wanxiang-agent.service
[Unit]
Description=Wanxiang Agent API
After=network.target

[Service]
WorkingDirectory=/www/wwwroot/t.agents.com/wanxiang
Environment=WANXIANG_ROOT=/www/wwwroot/t.agents.com/wanxiang
Environment=WANXIANG_HTTP_ADDR=127.0.0.1:8088
ExecStart=/www/wwwroot/t.agents.com/wanxiang/server/wanxiang
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
